# systemd/hytale-panel.service
[Unit]
Description=Hytale Panel
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=hytale-panel
Group=hytale
WorkingDirectory=/opt/hytale-panel
EnvironmentFile=/etc/hytale-panel/env
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=10
LimitNOFILE=65535
KillSignal=SIGTERM
TimeoutStopSec=15

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/hytale-panel/data /tmp /home/hytale/server
PrivateTmp=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false

[Install]
WantedBy=multi-user.target
